# src/kanberra-system-shutdown.service.in
[Unit]
Description=Play Shutdown Sound
DefaultDependencies=no
Before=shutdown.target

# A little optimization under the assumption that we are using the
# standard freedesktop theme
ConditionPathExists=/usr/share/sounds/freedesktop/stereo/system-shutdown.oga

[Service]
Type=oneshot
ExecStart=-@bindir@/kanberra-boot system-shutdown

[Install]
WantedBy=halt.target poweroff.target
Also=kanberra-system-shutdown-reboot.service kanberra-system-shutdown-reboot.service
